# .config/systemd/user/emacs.service
[Unit]
Description=Emacs: the extensible, self-documenting text editor

[Service]
Type=forking
ExecStart=/usr/bin/emacs --daemon=main
ExecStop=/usr/bin/emacsclient -s main --eval "(kill-emacs)"
Restart=always

[Install]
WantedBy=default.target
