# .config/autostart/compton.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=compton
Icon=compton
Exec=compton
TryExec=compton
Categories=Utility;